# worker/supervisord.conf
[supervisord]
nodaemon=true

[program:openvpn]
command=openvpn --config /root/worker.ovpn
